# inittab
::sysinit:/etc/init.d/rcS
::respawn:/usr/bin/runsvdir /var/service
